# main/draft/Testing/test.crontab
59 23 * * * python3 /app/upload.py >> /var/log/mytest.log 2>&1
